# docs/com.cbyb.prototype.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cbyb.prototype</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nathanmeyer/Dev/C-by-B-code/Prototype/venv/bin/python</string>
        <string>-m</string>
        <string>cbyb.app</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nathanmeyer/Dev/C-by-B-code/Prototype</string>

    <key>UserName</key>
    <string>nathanmeyer</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/nathanmeyer/Dev/C-by-B-code/Prototype/logs/server-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nathanmeyer/Dev/C-by-B-code/Prototype/logs/server-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
